# etc/systemd/user/openclaw-viz.service
[Unit]
Description=OpenClaw Shared Viz Service (mermaid/graphviz/chartjs render)
Documentation=https://github.com/DarojaAI/openclaw-gateway
# Soft dependency on the gateway (Wants=, not PartOf=).
#
# 2026-06-09 incident regression: this file previously had
#   Before=openclaw-gateway.service
#   PartOf=openclaw-gateway.service
# which caused a viz service reload to SIGTERM the gateway. PartOf= is the
# wrong directive for "I want to start with the gateway" — it cascades
# stop/restart lifecycle events to the parent unit. Wants= is correct.
#
# Do NOT reintroduce PartOf= or Before= here. If you need to coordinate
# start order, use After= (which is a start-order hint, not a lifecycle
# cascade). See AGENTS.md post-mortem 2026-06-09 for the full root cause.
Wants=openclaw-gateway.service
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/desktopuser/.openclaw/services/viz
ExecStart=/usr/bin/node /home/desktopuser/.openclaw/services/viz/render-server.js
Environment=VIZ_PORT=8766
Environment=VIZ_CACHE_DIR=/home/desktopuser/.openclaw/services/viz/cache
Environment=VIZ_HOST=localhost
Environment=NODE_COMPILE_CACHE=/var/tmp/openclaw-compile-cache
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=300
TimeoutStopSec=15
SuccessExitStatus=0 143
KillMode=control-group

[Install]
WantedBy=default.target
